# infrastructure/socket/supervisord.ini
[supervisord]
nodaemon=true
user=root

[program:httpd]
command = /usr/sbin/httpd -DFOREGROUND
autorestart = true
stderr_logfile = /var/www/html/storage/logs/apache.err.log
stdout_logfile = /var/www/html/storage/logs/apache.out.log

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/redis-worker.log


[program:reverb]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan reverb:start --port=8080
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/reverb.log
